# systemd/budget-agent-report.timer
[Unit]
Description=Budget Agent Monthly Report (last month)

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
