# systemd/budget-agent.timer
[Unit]
Description=Run Budget Agent weekly

[Timer]
OnCalendar=Sun *-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
